# iPaaS/project-root/internal/cloudapp_agent/doc/linux/serivces/ys-cloudapp.target
[Unit]
Description=yuansuan cloudapp target

[Install]
WantedBy=multi-user.target